# rootfs/etc/profile.d/03_prompt.sh
source /etc/profile.d/kube-ps1.sh
PS1='[\u@\h \W $(kube_ps1)]\$ '
